# systemd_units/rm-combined.service
[Unit]
Description=RM Combined Server UI proxy
After=network.target rm-mother.service

[Service]
Type=simple
User=joe
WorkingDirectory=/home/joe/sparky/stamp_engine
ExecStart=/usr/bin/python3 /home/joe/sparky/stamp_engine/combined_server.py
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=120
StartLimitBurst=3
StandardOutput=append:/home/joe/sparky/logs/combined.log
StandardError=append:/home/joe/sparky/logs/combined.log

[Install]
WantedBy=multi-user.target